# php/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-error.log

[program:websocket]
command=php /var/www/symfony_docker/bin/console app:websocket-server
autorestart=true
stdout_logfile=/var/log/websocket.log
stderr_logfile=/var/log/websocket-error.log
